# docker/crontab/crontab.sh
#!/bin/bash
ls

php artisan queue:work --queue=bizon_export